# deploy/rentalall.service
[Unit]
Description=RentalAll Django Application
After=network.target postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/rentalall/backend
Environment="PATH=/root/rentalall/backend/venv/bin"
ExecStart=/root/rentalall/backend/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8001 \
    --timeout 60 \
    --access-logfile /var/log/rentalall/access.log \
    --error-logfile /var/log/rentalall/error.log \
    rentalall.wsgi:application

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

